# doc/deploy/home/ubuntu/supervisor/webhooks.conf
[program:sisplan-github-webhooks]
directory=/home/ubuntu/webhooks
user=ubuntu

; command=python gitlab-webhooks.py
command=python webhooks.py

autorestart=true
