# doc/deploy/home/ubuntu/supervisor/webhooks.conf
[program:sisplan-github-webhooks]
directory=/media/datos/proyectos/elarteylatecnologia/mincultura/sgca/doc/deploy/home/ubuntu/webhooks
user=pardo
command=python webhooks.py
autorestart=true
